# deployment/backup/trade-show-app-backup.timer
[Unit]
Description=Run Trade Show App Backup Daily

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target


